# btrfs-backup.timer
[Unit]
Description=Weekly trigger for btrfs-backup.service

[Timer]
OnCalendar=Sun *-*-* 03:15:00
Persistent=true

[Install]
WantedBy=timers.target